# httpd/start_httpd.service
[Unit]
Description=Start instances httpd 

[Service]
Type=oneshot
ExecStart=/vagrant/httpd/start_httpd.sh

[Install]
WantedBy=multi-user.target
